# BackEnd/NewsCrawling/crontab
0 * * * * /usr/bin/python3 /usr/src/app/NewsCrawling.sh > /proc/1/fd/1 2>/proc/1/fd/2
